# internal/bootstrap/verta.service
# verta - enterprise secure mail server
#
# Runs as root: it binds privileged ports (25/465/587/143/993/110/995)
# and delivers to system users' Maildirs with their own UID/GID. The
# hardening below still confines it.

[Unit]
Description=verta - enterprise secure mail server
Documentation=https://github.com/ostap-mykhaylyak/verta
After=network-online.target
Wants=network-online.target
# Stop hammering if startup keeps failing (e.g. port 25 already in
# use): after 5 failures within 60s systemd gives up and the unit
# stays failed, so the problem is visible instead of a restart loop.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/sbin/verta start --config /etc/verta/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s
TimeoutStopSec=30s
KillSignal=SIGTERM

LimitNOFILE=65536

RuntimeDirectory=verta
RuntimeDirectoryMode=0750

# Hardening.
# ProtectSystem=full keeps /usr, /boot and /etc read-only while leaving
# the data hierarchy (/var, /home, /srv) writable — verta writes
# mailboxes wherever the domain files point them (/var/mail, a system
# user's home, a bind mount), so it cannot use the stricter
# ProtectSystem=strict without listing every mail location. For tighter
# confinement, set ProtectSystem=strict and add each mail path (plus
# /var/log/verta and /var/lib/verta) to ReadWritePaths.
ProtectSystem=full
# Maildir delivery to system users writes under /home.
ProtectHome=false
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
